# global-dotfiles/etc/udev/rules.d/83-gamingkeyboard.rules
ACTION=="add", KERNEL=="hidraw[0-9]*", SUBSYSTEM=="hidraw", SUBSYSTEMS=="usb", ATTRS{idVendor}=="04f3", ATTRS{idProduct}=="152e", RUN+="/home/ananthu/.config/scripts/cosmicbyte_keyboard.sh"
