# global-dotfiles/etc/systemd/system/script.service
[Unit]
Description=Startup_script

[Service]
ExecStart=/usr/local/bin/startup_script

[Install]
WantedBy=multi-user.target 
